# daemon/mail.mount.in
[Mount]
Type=mail
Exec=@libexecdir@/gvfsd-mail
AutoMount=true

